# tmp/supervisor_eth.conf
[program:eth]
command=/usr/bin/node /home/erc20-wallet/bin/eth.js
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/eth.err.log
stdout_logfile=/var/log/supervisor/eth.out.log